# root/etc/s6-overlay/s6-rc.d/flask/run
#!/command/execlineb -P
with-contenv
cd /opt/app
s6-setuidgid flask gunicorn --bind 0.0.0.0:8080 app:app
