# root/etc/s6-overlay/s6-rc.d/flask/run
#!/command/execlineb -P
with-contenv
export FLASK_APP /opt/app/app.py
s6-setuidgid flask flask run --host=0.0.0.0 --port=8080
